# contrib/voxgo.service
[Unit]
Description=voxgo dictation daemon
After=graphical-session.target

[Service]
ExecStart=%h/go/bin/voxgo daemon
Restart=on-failure

[Install]
WantedBy=default.target
